# contrib/super/apis.ini
[program:apis]
process_name=%(program_name)s
command=run_apis
directory=/opt/malwarescan
umask=002
autostart=true
autorestart=unexpected
exitcodes=0,2
user=malwarescan
stdout_logfile=/opt/malwarescan/logs/apis.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/opt/malwarescan/logs/apis.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
