# contrib/super/wsclient.ini
[program:wsclient]
process_name=%(program_name)s
command=run_wsclient
directory=/opt/malwarescan
umask=002
autostart=true
autorestart=unexpected
startsecs=10
startretries=5
exitcodes=0,2
user=malwarescan
stdout_logfile=/opt/malwarescan/logs/wsclient.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/opt/malwarescan/logs/wsclient.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
